# ui.service
[Unit]
After=network.target

[Service]
ExecStart=/bin/bash /etc/systemd/system/ui.sh

[Install]
WantedBy=default.target
